# systemd/jack-silence-detector.service.in
[Unit]
Description=jack-connection-manager
After=jack.service
Requires=jack.service

[Install]
# There is no multi-user.target for a user service, therefore we use default.target
WantedBy=default.target

[Service]
Type=simple
ExecStart=@bin_dir@/jack-silence-detector